# keetech/docker-odoo/odoo/files/supervisord.conf
[supervisord]
nodaemon=true                                   ; run supervisord as a daemon
loglevel=error                                  ; info, debug, warn, trace

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)
chmod=0777                 ; socked file mode (default 0700)
;chown=nobody:nogroup       ; socket file uid:gid owner
;username=user              ; (default is no username (open server))
;password=123               ; (default is no password (open server))

[program:odoo]
command=/opt/odoo/odoo-bin -c /opt/odoo/delta-config --db_host=%(ENV_PG_HOST)s --db_port=%(ENV_PG_PORT)s -r %(ENV_PG_USER)s -w %(ENV_PG_PASSWORD)s -u all
stdout_logfile=/var/log/odoo/odoo-update.log
redirect_stderr=true
redirect_stdout=true

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
